# setup/etc/cron.d/cld
@reboot root bash -lc "ip addr add 172.17.0.250/16 dev docker0"
*/5 * * * * root bash -lc "/var/cld/crons/hypervisor_parser ; /var/cld/crons/cldidgen"
