# setup/etc/cron.d/cld
*/5 * * * * root bash -lc "/var/cld/bin/cld-groupparser"
#* * * * * root bash -lc "/var/cld/bot/bin/access_list_to_ssh" &>/dev/null
#0 12 * * * root bash -lc "/var/cld/modules/etcbackup/bin/cld-etcbackup"