# setup/etc/cron.d/cld_etc_backup
0 12 * * * root bash -lc "/var/cld/crons/etc_backup"