# setup/etc/cron.d/cld_backup_checker
0 17 * * * root bash -lc "/var/cld/bot/bin/backup_checker send"